# supervisor/frpc.conf
[program:srv.frpc]
command=/opt/frp/bin/frpc -c /opt/frp/frpc.ini
directory=/opt/frp
user=root

logfile_maxbytes = 20MB
logfile_backups = 10

autorestart=true
startsecs=0
